# sw/pistackmond.service.in
[Unit]
Description=PiStackMon Daemon

[Service]
User=root
EnvironmentFile=/etc/default/pistackmond
ExecStart=PREFIX/bin/pistackmond -s -b ${BRIGHTNESS}

[Install]
WantedBy=local-fs.target
