# sw/pistackmond.service.in
[Unit]
Description=PiStackMon Daemon

[Service]
User=root

ExecStart=PREFIX/bin/pistackmond -s

[Install]
WantedBy=local-fs.target
