# sw/pistackmond.service
[Unit]
Description=PiStackMon Daemon

[Service]
User=root

ExecStart=/usr/local/sbin/pistackmond

[Install]
WantedBy=local-fs.target
